# addons/mindctrl/rootfs/etc/s6-overlay/s6-rc.d/traefik/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting traefik..."

exec /usr/bin/run_traefik.sh
