# addons/mindctrl/rootfs/etc/s6-overlay/s6-rc.d/multiserver/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting multiserver..."

exec /usr/bin/run_multiserver.sh
