# addons/mindctrl/rootfs/etc/s6-overlay/s6-rc.d/tracking/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting tracking server..."

exec /usr/bin/run_tracking.sh
