# mindctrl-addon/rootfs/etc/s6-overlay/s6-rc.d/dapr-dashboard/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting dashboard..."

exec /usr/bin/run_dashboard.sh
