# mlflow-addon/rootfs/etc/s6-overlay/s6-rc.d/gateway/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting gateway..."

exec /usr/bin/run_gateway.sh
